# meta-air-quality/recipes-apps/sensors/files/sensor-pms.service
[Unit]
Description=Sensor-pms5003 reading

[Service]
Type=oneshot
ExecStart=/usr/bin/sensor-pms /dev/ttyAMA0

[Install]
WantedBy=multi-user.target

