# meta-air-quality/recipes-apps/sensors/files/sensor-temp.service
[Unit]
Description=Sensor temperature reading

[Service]
ExecStart=/usr/bin/sensor-temp
Restart=oneshot

[Install]
WantedBy=multi-user.target
